# agents/macos/com.availcal.export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.availcal.export</string>

    <!-- install.sh rewrites this to the launcher inside the signed app bundle it
         builds (__INSTALL_DIR__/AvailCal.app/Contents/MacOS/availcal). Running
         the bundle's executable gives the job a stable macOS app identity, so
         the Calendar (TCC) prompt appears and the grant is listed/persisted
         under "AvailCal". The launcher spawns the venv python as a child, which
         inherits the app as its responsible process. -->
    <key>ProgramArguments</key>
    <array>
        <string>__APP_EXEC__</string>
    </array>

    <!-- Run hourly. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run shortly after load/login so the first feed appears promptly. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Upload target + auth, injected at install time (never hard-coded here).
         Cloudflare Worker: AVAILCAL_AGENT_SAS_URL is the upload URL
         https://availcal.<domain>/raw/<Label>.json and AVAILCAL_AGENT_TOKEN is
         the Worker's AGENT_TOKEN. Azure: the SAS URL alone (token left blank).
         launchd does NOT inherit your shell env, so both must live here. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>AVAILCAL_AGENT_SAS_URL</key>
        <string>__SAS_URL__</string>
        <key>AVAILCAL_AGENT_TOKEN</key>
        <string>__TOKEN__</string>
        <!-- Cloudflare Access service token, only if the Worker host is fronted
             by Access (Zero Trust). Left blank otherwise. -->
        <key>AVAILCAL_AGENT_CF_ACCESS_CLIENT_ID</key>
        <string>__CF_ACCESS_ID__</string>
        <key>AVAILCAL_AGENT_CF_ACCESS_CLIENT_SECRET</key>
        <string>__CF_ACCESS_SECRET__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/export.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/export.err.log</string>

    <!-- A non-zero exit (fail-loud) is recorded in export.err.log for alerting.
         ProcessType is Standard (not Background) so macOS does not deprioritize
         the job or suppress the first-run Calendar (TCC) permission prompt — the
         launchd-spawned python needs its OWN grant, separate from Terminal's. -->
    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
